# cronJob.cron
45 15 * * 1-5 /usr/bin/python3 /home/ubuntu/auric-fe/run.py
